# deploy/ultimate-kea-dashboard.service
[Unit]
Description=Ultimate Kea DHCP Dashboard (Plugin Runtime)
After=network.target

[Service]
Type=simple
User=kea
WorkingDirectory=/opt/ultimate-kea-dashboard
ExecStart=/usr/bin/python3 bin/ultimate-kea-dashboard-plugin
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
